# deploy/bootstrap/glycoview-appliance-bootstrap.service
[Unit]
Description=GlycoView appliance bootstrap
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/opt/glycoview/stack/.env
ExecStart=/opt/glycoview/bootstrap/bootstrap.sh

[Install]
WantedBy=multi-user.target
